# goldenway.sh
#! /bin/bash
#
### BEGIN INIT INFO
# Provides: Oh!Friends
# Default-Start: 5
### END INIT INFO

nohup /home/nvidia/goldenway/integrate &
